# .config/systemd/user/docker-compose.timer
[Timer]
OnStartupSec=1
# Every 15 minutes
OnCalendar=*:0/15
# Persistent=True

[Install]
# WantedBy=timers.target
WantedBy=default.user

