# .config/systemd/user/docker-compose@.timer
[Unit]

[Timer]
# Every 15 minutes
OnCalendar=*:0/15
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target

