# .config/systemd/user/docker-compose@.timer
[Unit]

[Timer]
OnStartupSec=1
# Every 15 minutes
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target

