# Templates/systemd/myservice.service
[Unit]
Description=

[Service]
Type=oneshot
ExecStart=

Install]
WantedBy=multi-user.target
